# deploy/vor-telegram.service
[Unit]
Description=VOR Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/root/vor
Environment=PATH=/root/vor/venv/bin
ExecStart=/root/vor/venv/bin/python telegram_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
